# rootfs/etc/services.d/php-fpm8/run
#!/command/with-contenv bash

echo "Starting php-fpm83"

exec /usr/sbin/php-fpm83